# app/config/deploy_files/smartkits.socket
[Unit]
Description=gunicorn socket

[Socket]
ListenStream=/run/smartkits.sock

[Install]
WantedBy=sockets.target